# mateix-watch.service
[Unit]
After=syslog.service

[Service]
ExecStart=/usr/bin/mateixWatch

[Install]
WantedBy=default.target
